# crontab
#run python script every hour
0 * * * * python /app/home_finder.py > /proc/1/fd/1 2>/proc/1/fd/2
