# crontab
#run python script every hour
SHELL=/bin/bash
30 * * * *  /app/cron_job.sh > /proc/1/fd/1 2>/proc/1/fd/2
